# teacher-registry.service
[Unit]
Description=Teacher Registry API Service
After=network.target

[Service]
Type=simple
User=teacher-registry
ExecStart=/usr/local/bin/teacher-registry
Restart=always
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target